# maven-app/3.3-jdk-7/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:app-init]
command=/usr/local/sbin/app-init
autorestart=false
