# php-app/php53/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/local/sbin/apache2

[program:app-init]
command=/usr/local/sbin/app-init
