# sphinx-app/1.4/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:app-init]
command=/usr/local/sbin/app-init
autorestart=false
startsecs=0
startretries=0
