# php-app/php70/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/local/sbin/apache2

[program:app-init]
command=/usr/local/sbin/app-init
autorestart=false
startsecs=0
startretries=0
