# symfony-app/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/local/sbin/symfony-app-apache2.sh

[program:app-init]
command=/usr/local/sbin/symfony-app-init.sh
