# symfony-app/php55/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/local/sbin/symfony-app-apache2

[program:app-init]
command=/usr/local/sbin/symfony-app-app-init.sh
